# static/etc/supervisor/conf.d/cloudflared.conf
[program:cloudflared]
# Admin service: Cloudflare tunnel. Conditional on CLOUDFLARE_TOKEN being set
# in the container environment (see /usr/local/bin/cloudflared-tunnel wrapper).
# autorestart=unexpected + exitcodes=0 means: a clean exit 0 (no token) is
# left EXITED/disabled, while a non-zero crash is restarted.
command=/usr/local/bin/cloudflared-tunnel
user=wall-e
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/wall-e/cloudflared.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/wall-e/cloudflared.err.log
stderr_logfile_maxbytes=0
