# static/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/wall-e/nginx.out.log
stderr_logfile=/var/log/wall-e/nginx.err.log
